# deploy/ultra-rpi.service
[Unit]
Description=Ultra RPi Controller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ultra3
WorkingDirectory=/home/ultra3/ultra-rpi
Environment=ULTRA_CONFIG=/etc/ultra/machine.yaml
ExecStart=/home/ultra3/ultra-rpi/.venv/bin/python -m ultra.app
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
